# backend/templates/systemd/ipa-backup.timer.j2
[Unit]
Description="IdM/FreeIPA backup timer"

[Timer]
OnCalendar={{ schedule }}
Persistent=true
Unit=ipa-backup.service

[Install]
WantedBy=timers.target
